# deploy/systemd/algobond-agent-watchdog.timer
[Unit]
Description=AlgoBond Sentinel Watchdog Timer

[Timer]
OnBootSec=60
OnUnitActiveSec=30
AccuracySec=5

[Install]
WantedBy=timers.target
